# deploy/iss-display.service
[Unit]
Description=ISS Tracker Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/iss-tracker
ExecStart=%h/iss-tracker/.venv/bin/python -m iss_display.app.main
EnvironmentFile=%h/iss-tracker/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
